# rovr.service
[Unit]
Description=Rovr - YouTube Downloader Bot for Nostr
After=network.target

[Service]
Type=simple
User=ethan
WorkingDirectory=/home/ethan/code/rovr
Environment="RUST_BACKTRACE=1"
Environment="PATH=/home/ethan/code/rovr/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/ethan/code/rovr/target/release/rovr
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 